# etc/init.d/install_ergo
#!/bin/sh
### BEGIN INIT INFO
# Provides: install_ergo
# Required-Start:
# Required-Stop:
# Default-Start: 3 
# Default-Stop:
# Short-Description: grabs keys and installs ergo
# Description:
### END INIT INFO

do_start
{
PASSWORD="a_password"
UNIT_ID="a_id" 

#sets up apt sources
wget -O -  http://packages.ergotelescope.org/ergo.gpg.key | apt-key add - 

#sets up ergo and unattended-upgrades
apt-get update;
apt-get upgrade -y;
apt-get install -y ergo-telescope unattended-upgrades apt-listchanges;
apt-get update;

#changes the password for the unit with the format password_unit id 
# eg password_555
echo -e "$PASSWORD_$UNIT_ID\n$PASSWORD_$UNIT_ID" |(passwd pi)
echo "$UNIT_ID" > /etc/ERGO/unit_id

#adds a key for admins. 
mkdir -p /home/pi/.ssh
cat /root/master_key >> /home/pi/.ssh/authorized_keys

chmod 600 /etc/init.d/install_ergo
update-rc.d install_ergo remove
}
